# hamtrek@.service
[Unit]
Description=Hamtrek Server

[Service]
User=prad
Group=packetradio
ExecStart=/opt/hamtrek/hamtrek -h
StandardInput=socket
