# hamtrek@.service
[Unit]
Description=Hamtrek Server

[Service]
ExecStart=/opt/hamtrek/hamtrek -h
StandardInput=socket
